# nodejs/Add/supervisor.nodejs.conf
[program:nodejs]
command=/usr/bin/nodejs /app/node/default/start.js